# redis_docker/app/supervisor.conf
[supervisord]
nodaemon=true

[program:celery_prog]
command=celery -A tasks1 worker --loglevel=INFO



[program:flask_prog]
command=python3 app.py

